# services/mm-os-update.service
[Unit]
Description=MagicMirror OS Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/magicmirror-setup/scripts/update-os.sh
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
